# systemd/hoval-exporter.service
[Unit]
Description=HoxPi Prometheus Exporter (Modbus -> :9101/metrics)
After=hoval-bridge.service

[Service]
ExecStart=/usr/bin/python3 /opt/hoxpi/hoval_exporter.py
Restart=always
RestartSec=5
DynamicUser=yes

[Install]
WantedBy=multi-user.target
